# infra/host/combo-dev/combo-dev-storage-guard.timer
[Unit]
Description=Continuously check combo-dev bounded storage headroom

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=combo-dev-storage-guard.service

[Install]
WantedBy=timers.target
